# Dockerfiles/files/supervisord-rtl2832u.conf
[program:nodeBooter]
autostart=true
redirect_stderr=true
command=/bin/bash -lc "nodeBooter -d /nodes/%(ENV_NODENAME)s/DeviceManager.dcd.xml"
priority=10
